# roles/realm_join/files/httpd-gssproxy.service
.include /lib/systemd/system/httpd.service
[Service]
Environment=GSS_USE_PROXY=1
